# etc/systemd/ec-boot-api.service
[Unit]
Description=Exordos Core Boot API
After=network-online.target

[Service]
TimeoutStopSec=5
Restart=always
RestartSec=5s
KillSignal=SIGINT
ExecStart=ec-boot-api --config-file /etc/exordos_core/exordos_core.conf

[Install]
WantedBy=multi-user.target
